# infra/mac/com.myroproductions.stt-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myroproductions.stt-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=no</string>
        <string>-i</string>
        <string>/Users/myro-pro/.ssh/id_ed25519</string>
        <string>-L</string>
        <string>8765:localhost:8765</string>
        <string>linux-desktop</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>RestartThrottle</key>
    <integer>10</integer>
    <key>StandardErrorPath</key>
    <string>/tmp/stt-tunnel.log</string>
</dict>
</plist>
